# etc/systemd/system/letsencrypt.service
[Unit]
Description=Let's Encrypt SSL Certificate Renewal

[Service]
Type=simple
ExecStart=/usr/local/bin/letsencrypt-renew
